# python/crontab
* * * * * root echo "Hello $(date)" >> /var/log/cron.log 2>&1
* * * * * root /root/tmp.sh >> /var/log/cron.log 2>&1
* * * * * root python3 /root/main.py >> /var/log/cron.log 2>&1
